# build/docker/monaco-go_go-langserver/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D -e

[program:go-langserver]
command=langserver-go -mode ws -addr :4389 -trace -logfile /Users/mbana/monaco-go/go-langserver.log